# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

; 1. Qdrant vector database
[program:qdrant]
command=/usr/local/bin/qdrant --config-path /tmp/qdrant_config/config.yaml
directory=/app
autostart=true
autorestart=true
stdout_logfile=/tmp/qdrant.log
stderr_logfile=/tmp/qdrant.log
priority=10

; 2. PageIndex microservice (RAG segmentation + retrieval)
[program:pageindex]
command=uvicorn pageindex.service:app --host 0.0.0.0 --port 8001 --workers 1
directory=/app
environment=PYTHONPATH="/app"
autostart=true
autorestart=true
stdout_logfile=/tmp/pageindex.log
stderr_logfile=/tmp/pageindex.log
priority=20
startsecs=8

; 3. Main FastAPI
[program:api]
command=uvicorn api.main:app --host 0.0.0.0 --port 8000 --workers 1
directory=/app
environment=PYTHONPATH="/app"
autostart=true
autorestart=true
stdout_logfile=/tmp/api.log
stderr_logfile=/tmp/api.log
priority=30
startsecs=12

; 4. Streamlit dashboard (external port 7860)
[program:dashboard]
command=streamlit run dashboard/app.py --server.port 7860 --server.address 0.0.0.0 --server.headless true --server.fileWatcherType none
directory=/app
environment=PYTHONPATH="/app",API_BASE_URL="http://localhost:8000"
autostart=true
autorestart=true
stdout_logfile=/tmp/dashboard.log
stderr_logfile=/tmp/dashboard.log
priority=40
startsecs=20
